# deploy/launchd/com.ctk03272.mcp-korean-stock.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.ctk03272.mcp-korean-stock</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>__RUNNER_SCRIPT__</string>
      <string>__APP_ROOT__</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>__APP_ROOT__/current</string>
    <key>StandardOutPath</key>
    <string>__APP_ROOT__/shared/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__APP_ROOT__/shared/logs/stderr.log</string>
  </dict>
</plist>
